# tpm4/config/includes.chroot/etc/xdg/autostart/inotify-monitor.desktop
[Desktop Entry]
name=inotify-monitor
Type=Application
Terminal=true
Exec=/home/usuario/inotify-monitor.sh
Icon=mc
Dscription=TPM Monitor
