# tpm4/config/includes.chroot/etc/xdg/autostart/tpm-monitor.desktop
[Desktop Entry]
name=tpm-monitor
Type=Application
Terminal=true
Exec=gjs /home/usuario/tpm-monitor.js
Icon=mc
Dscription=TPM Monitor
